# .config/crontabs/user_crontab
@weekly				trash-empty 14
@weekly				yay -Sc
*/10 * * * *        $HOME/.local/bin/systemd-failure-notification.py
@weekly				find $HOME/Pictures/Screenshots -mtime +30 -exec rm {} \;

